# algo_api.service
[Unit]
Description=Algo service

[Service]
User=bn
WorkingDirectory=/home/bn/algo_testing
ExecStart=/usr/bin/python3.8 /home/bn/algo_testing/algo.py
Restart=always

[Install]
WantedBy=multi-user.target
